# scripts/uncivserver.service
[Unit]
Description=UncivServer.xyz
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=bun start
Restart=always
RestartSec=1
StartLimitBurst=3600
User=nobody
Group=nogroup
# Replace `/bun/directory` with `which bun` output
ExecSearchPath=`/bun/directory`
# replace `/path/to/home/UncivServer.xyz` with server repo path
WorkingDirectory=/path/to/home/UncivServer.xyz

[Install]
WantedBy=multi-user.target
