# scripts/uncivserver.service
[Unit]
Description=UncivServer.xyz
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=bun start
Restart=always
RestartSec=1
User=admons
# Replace `/path/to/bun/directory/bin` with `dirname $(which bun)` output
ExecSearchPath=/path/to/bun/directory/bin
# replace `/path/to/UncivServer.xyz` with server repo path
WorkingDirectory=/path/to/UncivServer.xyz

[Install]
WantedBy=multi-user.target
